# recipes-panda/led-daemon/files/led-daemon.service
[Unit]
Description=led-daemon service

[Service]
ExecStart=/usr/bin/led-daemon

[Install]
WantedBy=multi-user.target
